# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:rabbitmq_app]
command=bash -c "python3 file_gen_service/rabbit_folder/consumer.py"
autostart=true
autorestart=true
startsecs=10
numprocs=5
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=supervisor.log
stdout_syslog=true
stderr_logfile=supervisor.log
stderr_syslog=true
